# resolver/deploy/humanity-resolver.service
[Unit]
Description=Humanity Commons Resolver (v0.5)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/resolver
EnvironmentFile=/opt/resolver/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
